# project-tools/mcp-hooks-server/com.mcp-augment.plist
<!-- To install: cp this file to ~/Library/LaunchAgents/ && launchctl load ~/Library/LaunchAgents/com.mcp-augment.plist -->
<!-- To uninstall: launchctl unload ~/Library/LaunchAgents/com.mcp-augment.plist && rm ~/Library/LaunchAgents/com.mcp-augment.plist -->
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mcp-augment</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PROJECT_DIR__/project-tools/mcp-hooks-server/start-servers.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/mcp-augment-launchd.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/mcp-augment-launchd.stderr.log</string>
</dict>
</plist>
